# Sources/KeyPathHelper/launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service identifier -->
    <key>Label</key>
    <string>com.keypath.helper</string>

    <!-- Mach service for XPC communication -->
    <key>MachServices</key>
    <dict>
        <key>com.keypath.helper</key>
        <true/>
    </dict>

    <!-- Run as root -->
    <key>UserName</key>
    <string>root</string>

    <!-- Program arguments -->
    <key>ProgramArguments</key>
    <array>
        <string>/Library/PrivilegedHelperTools/com.keypath.helper</string>
    </array>

    <!-- On-demand activation via XPC (not always resident) -->
    <key>KeepAlive</key>
    <false/>

    <!-- Don't run at system boot (XPC will launch on demand) -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Standard output/error logging -->
    <key>StandardOutPath</key>
    <string>/var/log/com.keypath.helper.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/com.keypath.helper.stderr.log</string>
</dict>
</plist>
