# Sources/KeyPathApp/com.keypath.kanata.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.keypath.kanata</string>
    
    <!-- SMAppService: program path relative to the app bundle root -->
    <key>BundleProgram</key>
    <string>Contents/Library/KeyPath/kanata-launcher</string>
    
    <!-- Program arguments: program path (required even with BundleProgram) and config -->
    <!-- Note: Config path uses ~ which will be expanded by launchd -->
    <!-- TCP port defaults to 37001, can be overridden via environment if needed -->
    <key>ProgramArguments</key>
    <array>
        <string>Contents/Library/KeyPath/kanata-launcher</string>
        <string>--port</string>
        <string>37001</string>
        <string>--debug</string>
        <string>--log-layer-changes</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/var/log/com.keypath.kanata.stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/var/log/com.keypath.kanata.stderr.log</string>
    
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>256</integer>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <key>AssociatedBundleIdentifiers</key>
    <array>
        <string>com.keypath.KeyPath</string>
    </array>
</dict>
</plist>
